# luci-app-netopt/root/etc/init.d/netopt
#!/bin/sh /etc/rc.common

START=99          # 启动优先级，99 表示网络启动后运行
STOP=10           # 停止优先级（可选）

start() {
    config_load netopt
    config_get enable config enable
    if [ "$enable" = "1" ]; then
        config_get eee config eee '0'
        config_get rpsxps config rpsxps '0'
        config_get interfaces config interfaces

        # 处理 interfaces 为 list 类型，拼接成空格分隔字符串
        local iface_list=""
        for iface in $interfaces; do
            iface_list="$iface_list $iface"
        done
        iface_list="${iface_list# }"  # 去掉开头空格

        logger -t netopt "Boot optimization enabled. Applying settings..."

        # 调用 rpcd 后端方法（你的 shell 脚本）
        ubus call netopt apply_optimization "{\"eee\":\"$eee\",\"rpsxps\":\"$rpsxps\",\"interfaces\":\"$iface_list\"}"
        if [ $? -eq 0 ]; then
            logger -t netopt "Boot optimization applied successfully"
        else
            logger -t netopt "Boot optimization failed"
        fi
    else
        logger -t netopt "Boot optimization disabled (enable=0)"
    fi
}

stop() {
    # 可选：停止时不做任何事，或者清理临时状态
    logger -t netopt "Netopt service stopped"
}

boot() {
    start
}
